# DEBIAN/postinst
systemctl enable blankwatch.service
systemctl start blankwatch.service