# DEBIAN/postinst
systemctl enable blankwatch.service
systemctl stat blankwatch.service