# etc/systemd/blank_slot.service
# /etc/systemd/system/blankwatch.service
[Unit]
Description=BlankWatch
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/scripts/blankwatch.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
